# distribution/systemd/nostdb.service
[Unit]
Description=NostDB single-node database daemon
Documentation=https://github.com/nostdb/nostdb-server
After=network.target

[Service]
Type=simple
User=nostdb
Group=nostdb
ExecStart=/usr/local/bin/nostd serve --config /etc/nostdb/server.toml
Restart=on-failure
RestartSec=2s
TimeoutStopSec=10s
StateDirectory=nostdb
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/nostdb
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
